# gunicorn/gunicorn.socket
[Unit]
Description=gunicorn socket
[Socket]
ListenStream=/run/llmchatbot-gunicorn.sock
[Install]
WantedBy=sockets.target